# com.phfactor.fa-web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.phfactor.fa-web</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/pfh/code/fa-web/.build/release/App</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/pfh/code/fa-web</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/pfh/code/fa-web/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/pfh/code/fa-web/logs/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
